# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn dgrenier.wsgi:application -c /opt/gunicorn.py

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log